# roles/device/firewall/files/etc/network/if-pre-up.d/ipv4-firewall-rules
#!/bin/sh

echo "Loading IPv4 firewall rules"
iptables-restore  < /etc/netfilter/iptables.txt
