# roles/device/router/files/etc/network/if-pre-up.d/ipv6-firewall-rules
#!/bin/sh

echo "Loading IPv6 firewall rules"
ip6tables-restore < /etc/netfilter/ip6tables.txt
